# services/odoo_server/hw_proxy.service
[Unit]
Description=HW Proxy FastAPI Service
After=network.target

[Service]
User=hw_user
Group=hw_user
WorkingDirectory=/opt/hw_proxy
Environment="PATH=/opt/hw_proxy/.venv/bin"
ExecStart=/opt/hw_proxy/.venv/bin/uvicorn hw_proxy.main:app --host 0.0.0.0 --port 9002 --workers 2 --log-level info
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hw_proxy

# Security hardening (recommended)
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target